# ops/deploy/supervisor/glintstone-staging-web.ini
; Glintstone staging web app — uvicorn, port 8004, separate deploy dir.
; Single worker (staging is low-traffic). autostart=false: staging is brought
; up on demand. Created by provision-staging.sh; mirrored here so deploy.sh's
; sync of ops/deploy/supervisor/ doesn't delete it from /etc/supervisor.d/.

[program:glintstone-staging-web]
command=/var/www/glintstone-staging/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8004
directory=/var/www/glintstone-staging
user=deploy
autostart=false
autorestart=true
environment=APP_ENV="staging"
stderr_logfile=/var/log/glintstone-staging-web.err.log
stdout_logfile=/var/log/glintstone-staging-web.out.log
