# ops/deploy/supervisor/glintstone-mcp.ini
; Glintstone MCP HTTP server.
;
; mcp.server_http is not implemented yet (mcp/transport/ is empty). The stanza
; is registered so the service exists, but autostart=false keeps supervisord
; from crash-looping it on boot. Flip autostart back to true once
; mcp.server_http lands.

[program:glintstone-mcp]
command=/var/www/glintstone/current/venv/bin/python -m mcp.server_http --host 127.0.0.1 --port 8005
directory=/var/www/glintstone/current
user=deploy
autostart=false
autorestart=true
stderr_logfile=/var/log/glintstone-mcp.err.log
stdout_logfile=/var/log/glintstone-mcp.out.log
