# ops/deploy/supervisor/glintstone-web.ini
; Glintstone production web app — uvicorn (FastAPI, server-rendered Jinja2).
;
; Worker / connection-pool sizing: see glintstone-api.ini for the full
; 4-pools × DB_POOL_MAX math. api + web each run --workers 2; DB_POOL_MAX=10
; → 40 peak DB connections, well under the ~97 usable. The deploy health-gate
; alerts at ~70/97 before exhaustion (#472).
;
; This file is the AUTHORITATIVE copy: deploy.sh rsyncs ops/deploy/supervisor/
; to /etc/supervisor.d/ on every deploy and runs `supervisorctl update` (#472).
; APP_VERSION is omitted on purpose — version.txt (written by deploy.sh) is the
; source of truth for the release tag (core.version.release_tag()).

[program:glintstone-web]
command=/var/www/glintstone/current/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8002 --workers 2
directory=/var/www/glintstone/current
user=deploy
autostart=true
autorestart=true
stderr_logfile=/var/log/glintstone-web.err.log
stdout_logfile=/var/log/glintstone-web.out.log
environment=HOME="/home/deploy",APP_ENV="production",DB_POOL_MAX="10"
