# ops/deploy/supervisor/glintstone-api.ini
; Glintstone production API — uvicorn (FastAPI).
;
; Worker / connection-pool sizing (production, Postgres max_connections=100):
;   - api and web each run uvicorn with --workers 2 → 4 worker processes total.
;   - Each worker opens ONE psycopg pool (core.database.init_pool), so peak DB
;     connections from the app = 4 pools × DB_POOL_MAX.
;   - DB_POOL_MAX=10 → 4 × 10 = 40 peak, well under the ~97 usable connections
;     (100 − 3 superuser_reserved). Do NOT raise workers or DB_POOL_MAX without
;     re-checking this math — the deploy health-gate now alerts at ~70/97 (#472).
;
; This file is the AUTHORITATIVE copy: deploy.sh rsyncs ops/deploy/supervisor/
; to /etc/supervisor.d/ on every deploy and runs `supervisorctl update`, so the
; repo — not a manual VPS edit — owns the worker count and env (#472, #444).
;
; APP_VERSION is intentionally NOT set here: core.version.release_tag() reads the
; per-release version.txt that deploy.sh writes (version.txt wins over the env
; var), so a value baked into this file would only go stale. APP_ENV is set both
; here and in shared/.env; either alone is sufficient.

[program:glintstone-api]
command=/var/www/glintstone/current/venv/bin/uvicorn api.main:app --host 127.0.0.1 --port 8001 --workers 2
directory=/var/www/glintstone/current
user=deploy
autostart=true
autorestart=true
stderr_logfile=/var/log/glintstone-api.err.log
stdout_logfile=/var/log/glintstone-api.out.log
environment=HOME="/home/deploy",APP_ENV="production",DB_POOL_MAX="10"
