# ops/deploy/supervisor/glintstone-staging-api.ini
; Glintstone staging API — uvicorn, port 8003, separate deploy dir.
; Single worker (staging is low-traffic). autostart=false: staging is brought
; up on demand. Created by provision-staging.sh; mirrored here so deploy.sh's
; sync of ops/deploy/supervisor/ doesn't delete it from /etc/supervisor.d/.

[program:glintstone-staging-api]
command=/var/www/glintstone-staging/venv/bin/uvicorn api.main:app --host 127.0.0.1 --port 8003
directory=/var/www/glintstone-staging
user=deploy
autostart=false
autorestart=true
environment=APP_ENV="staging"
stderr_logfile=/var/log/glintstone-staging-api.err.log
stdout_logfile=/var/log/glintstone-staging-api.out.log
