# ops/deploy/supervisor/glintstone-crawler.ini
; Glintstone CDLI image crawler — long-running background process.
;
; Walks the artifacts table, fetches each CDLI artifact page at the mandatory
; 60s/req rate limit, parses out the image manifest, downloads originals,
; generates WebP thumbnails, and uploads everything to R2. Resumable via an
; on-disk checkpoint, so respawn-on-crash is safe.
;
; Expected runtime at steady state: ~245 days for the full 350k-artifact
; backfill. Auto-restart absorbs transient failures (network blips, parser
; misses, CDLI 5xx) without losing progress.
;
; .env is picked up via pydantic-settings reading from the working directory
; (the release dir's .env is a symlink to /var/www/glintstone/shared/.env,
; matching how api/app get their config).

[program:glintstone-crawler]
command=/var/www/glintstone/current/venv/bin/python -m ops.scripts.cdli_image_crawler
directory=/var/www/glintstone/current
user=deploy
autostart=true
autorestart=true
startsecs=30
startretries=10
stopsignal=TERM
stopwaitsecs=120
stderr_logfile=/var/log/glintstone-crawler.err.log
stdout_logfile=/var/log/glintstone-crawler.out.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
environment=HOME="/home/deploy"
